# example.home-blink.service
[Unit]
Description=Home Blink Python Service
After=network.target

[Service]
WorkingDirectory=/home/{USERNAME}/Projects/home-blink
ExecStart=/home/{USERNAME}/Projects/home-blink/venv/bin/python -m home_blink.main
Restart=always
RestartSec=5

StandardOutput=append:/home/{USERNAME}/Projects/home-blink/log/home-blink.og
StandardError=append:/home/{USERNAME}/Projects/home-blink/log/home-blink.log

[Install]
WantedBy=multi-user.target
